# alignak/systemd/alignak-scheduler.service
[Unit]
Description=Alignak Scheduler
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/alignak-scheduler -d -c /etc/alignak/daemons/schedulerd.ini

[Install]
WantedBy=multi-user.target
